# greta-cloudrun/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Express orchestrator (runs initializeProject, manages MongoDB, GCS) ───────
[program:server]
command=node /app/server.js
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Vite dev server (started by server.js after project init completes) ───────
[program:frontend]
command=bun run dev -- --host 0.0.0.0 --port %(ENV_VITE_PORT)s
directory=/app/project/frontend
autostart=false
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/frontend.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/frontend.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

; ── Python FastAPI backend (started by server.js after project init completes)
[program:backend]
command=uvicorn server:app --host 0.0.0.0 --port %(ENV_BACKEND_PORT)s --reload
directory=/app/project/backend
autostart=false
autorestart=true
startretries=3
environment=MONGO_URL="mongodb://localhost:%(ENV_MONGO_PORT)s",PYTHONPATH="/app/project/backend"
stdout_logfile=/var/log/supervisor/backend.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/backend.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
